# deployment_scripts/puppet/modules/heka/templates/hekad.upstart.conf.erb
# <%= @_service_name %>

description     "<%= @_service_name %>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

<% unless @_run_as_root -%>
pre-start script
    touch <%= @log_file %>
    chown <%= @heka_user %>:<%= @heka_user %> <%= @log_file %>
end script
<% end -%>

script
    # https://bugs.launchpad.net/lma-toolchain/+bug/1543289
    ulimit -n 102400
    exec start-stop-daemon --start <%= @_run_as_root ? "" : " --chuid #{ @heka_user }" %> --exec <%= @hekad_wrapper %> 2>><%= @log_file %>
end script
