# deployment_scripts/puppet/modules/heka/templates/logrotate.cron.erb
#!/bin/bash
# managed by puppet

test -x /usr/sbin/logrotate || exit 0

LOCK_FILE="/var/lock/logrotate.<%= @service_name %>.lock"

lock() {
    exec 903>$LOCK_FILE
    flock -n 903 && return 0 || return 1
}

unlock() {
    flock -u 903
}

fail() {
    if [ -z "$1" ]
    then
        MESSAGE="WARNING logrotate failed, no reason provided"
    else
        MESSAGE=$1
    fi
    /usr/bin/logger -t logrotate "<%= @service_name %> ${MESSAGE}"
    unlock
    exit 1
}

lock || fail "WARNING <%= @service_name %> logrotate flock failed, exiting"


TMP_FILE=$(/bin/mktemp)
nice ionice -c3 /usr/sbin/logrotate -s /var/lib/logrotate/<%= @service_name %>.status <%= @logrotate_conf %> >& $TMP_FILE
EXITVALUE=$?

if [ -f /etc/redhat-release ] || [ -f /etc/centos-release ];
then
    # Due to bug in logrotate on centos/rhel, it always returns 0. Use grep for
    # detect errors; exit code 1 is considered a success as no errors were
    # found.
    grep -q error $TMP_FILE
    EXITVALUE=$?
    EXPECTEDVALUE=1
else
    EXPECTEDVALUE=0
fi
rm $TMP_FILE

if [ "${EXITVALUE}" != "${EXPECTEDVALUE}" ]; then
    fail "ALERT exited abnormally with [${EXITVALUE}] (${EXPECTEDVALUE} was expected)"
fi

unlock
exit 0
